# share/muso.service
[Unit]
Description=Muso watcher mode

[Service]
Type=simple
ExecStart=/usr/bin/env muso watch

[Install]
WantedBy=default.target
